# sys-cluster/minio-bin/files/minio_initd
#!/sbin/openrc-run
# Copyright 2024 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="minio daemon"
description=""
command=/usr/bin/minio
command_args="${minio_args}"

supervisor="supervise-daemon"
command_background=yes

pidfile="/run/$RC_SVCNAME.pid"
output_log="$logfile"
error_log="$logfile"


depend() {
	need dmcrypt net
}
