# sys-cluster/seaweedfs/files/seaweedfs_initd
#!/sbin/openrc-run
# Copyright 2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

INSTANCE="${SVCNAME##*.}"

if [ -z "$INSTANCE" ] 
then
	exit 1
else
	seaweedfs_config="$INSTANCE"
fi

name="${SVCNAME} daemon"
description=""
command=/usr/bin/weed
command_args="${seaweedfs_args}"


#: "${SEAWEEDFS_USER:="seaweedfs"}"
#: "${SEAWEEDFS_GROUP:="seaweedfs"}"
: "${LOG_STDOUT:="/var/log/${SVCNAME}.stdout.log"}"
: "${LOG_STDERR:="/var/log/${SVCNAME}.stderr.log"}"

#command_user="$SEAWEEDFS_USER:$SEAWEEDFS_GROUP"
supervise_daemon_args="--stdout $LOG_STDOUT --stderr $LOG_STDERR"

depend() {
	need logger net
}

supervisor="supervise-daemon"
pidfile="/run/${INSTANCE}.pid"

#supervisor="start-stop-daemon"
#output_logger="/usr/bin/logger --tag lpmd"
#error_logger="/usr/bin/logger --tag lpmd_err"
#command_background=true

#start_pre(){
#	touch $LOG_STDOUT $LOG_STDERR
#	chown $SEAWEEDFS_USER:$SEAWEEDFS_GROUP $LOG_STDOUT $LOG_STDERR
#}

#start_pre(){
#	ebegin "$INSTANCE"
#}
