# net-misc/wsdd2/files/wsdd2_initd
#!/sbin/openrc-run
# Copyright 2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="wsdd2 daemon"
description=""
command=/usr/sbin/wsdd2
command_args=""

#supervise_daemon_args="--stdout $LOG_STDOUT --stderr $LOG_STDERR"

output_logger="/usr/bin/logger --tag wsdd2"
error_logger="/usr/bin/logger --tag wsdd2_err"

depend() {
	need logger net
}

supervisor="supervise-daemon"
pidfile="/run/wsdd2.pid"
