# sys-block/usb-nvme-adapter-discard/files/08-jmicron-usb-nvme-adapter-discard.rules
ACTION=="add|change", ATTRS{idVendor}=="152d", ATTRS{idProduct}=="0583", SUBSYSTEM=="scsi_disk" , ATTR{provisioning_mode}="unmap", RUN+="udev_delayed_unmap.sh"
